# infrastructure/services/php/symfony6/messenger-worker-dev.ini
[program:messenger-consume_ac_news_open_list]
command=bash -c "sleep 20 && php -d memory_limit=-1 /var/www/html/bin/console messenger:consume collector_ac_news_open_list"
user=root
numprocs=1
startsecs=0
autostart=true
autorestart=true
startretries=5
process_name=%(program_name)s_%(process_num)01d

[program:messenger-consume_ac_news_un_open_list]
command=bash -c "sleep 20 && php -d memory_limit=-1 /var/www/html/bin/console messenger:consume collector_ac_news_unopen_list"
user=root
numprocs=1
startsecs=0
autostart=true
autorestart=true
startretries=5
process_name=%(program_name)s_%(process_num)01d

[program:messenger-consume_ac_contact_import]
command=bash -c "sleep 20 && php -d memory_limit=-1 /var/www/html/bin/console messenger:consume collector_ac_contact_import"
user=root
numprocs=1
startsecs=0
autostart=true
autorestart=true
startretries=5
process_name=%(program_name)s_%(process_num)01d

[program:messenger-consume_ac_link_import]
command=bash -c "sleep 20 && php -d memory_limit=-1 /var/www/html/bin/console messenger:consume collector_ac_link_import"
user=root
numprocs=1
startsecs=0
autostart=true
autorestart=true
startretries=5
process_name=%(program_name)s_%(process_num)01d


[program:messenger-consume_sync-create]
command=bash -c "sleep 20 && php -d memory_limit=-1 /var/www/html/bin/console messenger:consume sync_create"
user=root
numprocs=1
startsecs=0
autostart=true
autorestart=true
startretries=5
process_name=%(program_name)s_%(process_num)01d

[program:messenger-consume_sync-validate]
command=bash -c "sleep 20 && php -d memory_limit=-1 /var/www/html/bin/console messenger:consume sync_validate"
user=root
numprocs=1
startsecs=0
autostart=true
autorestart=true
startretries=5
process_name=%(program_name)s_%(process_num)01d

[program:messenger-consume_sync-upload]
command=bash -c "sleep 20 && php -d memory_limit=-1 /var/www/html/bin/console messenger:consume sync_upload"
user=root
numprocs=1
startsecs=0
autostart=true
autorestart=true
startretries=5
process_name=%(program_name)s_%(process_num)01d

[program:messenger-consume_vimeo_import]
command=bash -c "sleep 20 && php -d memory_limit=-1 /var/www/html/bin/console messenger:consume vimeo_import"
user=root
numprocs=1
startsecs=0
autostart=true
autorestart=true
startretries=5
process_name=%(program_name)s_%(process_num)01d